# package/hong/luci-theme-advancedtomato/root/etc/uci-defaults/30_luci-theme-advancedtomato
#!/bin/sh
uci batch <<-EOF
	set luci.themes.advancedtomato=/luci-static/advancedtomato
	set luci.main.mediaurlbase=/luci-static/advancedtomato
	commit luci
EOF
exit 0
